# ec2-snapshot/systemd/sygnif-brain-context.service.d/harden.conf
[Service]
Environment=SYGNIF_BRAIN_CONTEXT_INTERVAL=180
