# ec2-snapshot/systemd/sygnif-trade-nl-publisher.service.d/harden.conf
[Service]
Environment=SYGNIF_NL_INTERVAL=90
