# ec2-snapshot/systemd/sygnif-dashboard-btc-terminal.service.d/10-port.conf
[Service]
Environment=SYGNIF_BTC_TERMINAL_PORT=8891
