# ec2-snapshot/systemd/sygnif-bybit-nl-feed.service.d/disable-hivemind.conf
[Service]
Environment=BYBIT_NL_HIVEMIND_DISABLED=1
