# ec2-snapshot/systemd/sygnif-bybit-nl-feed.service.d/harden.conf
[Service]
Environment=BYBIT_NL_POST_TIMEOUT_SEC=120
Environment=BYBIT_NL_INTERVAL=60
