# ec2-snapshot/systemd/sygnif-bybit-nl-feed.service.d/swarm-full.conf
[Service]
Environment=BYBIT_NL_SWARM_FULL=1
