# ec2-snapshot/systemd/sygnif-fast-reactor.service.d/momentum-tighten.conf
# Tighten the momentum trigger thresholds for current low-vol BTC regime.
#
# Why: with defaults TRIGGER_MOMENTUM_PCT=0.4 and TRIGGER_MOMENTUM_VOLX=1.5
# the momentum trigger fired ZERO times in a 7-day BTCUSDT 1m backtest run
# 2026-05-13. BTC has been drifting slowly in $79.5k-$82.5k, so individual
# 1m bars rarely move 0.4% AND have 1.5× volume spike.
#
# Lowered to 0.2% / 1.2× to recover the signal source. The intel-veto path
# inside fire_trade() still acts as a safety net against bad direction.
#
# Companion change: the new fib_sfp_confluence trigger (see
# experiments/fast_reactor_v2/) covers swing-failure-at-fib edges that
# momentum doesn't catch.

[Service]
Environment=SYGNIF_FAST_MOMENTUM_PCT=0.2
Environment=SYGNIF_FAST_MOMENTUM_VOLX=1.2
