# systemd/btc-predict-runner.service
[Unit]
Description=BTC next-bar ML predict (btc_predict_runner; flock with ruleprediction + finetune)
Documentation=file:///home/ubuntu/SYGNIF/scripts/btc_predict_runner_host.sh
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
EnvironmentFile=-/home/ubuntu/SYGNIF/.env
# Reads/writes under repo; skips quietly if SYGNIF_ML_LOCK_FILE is held.
ExecStart=/bin/bash /home/ubuntu/SYGNIF/scripts/btc_predict_runner_host.sh
StandardOutput=journal
StandardError=journal
