# systemd/bybit-stream-monitor.service
[Unit]
Description=Bybit WebSocket monitor (public orderbook/trades + demo private execution stream)
Documentation=file:///home/ubuntu/SYGNIF/scripts/bybit_stream_monitor.py
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/SYGNIF
EnvironmentFile=-/home/ubuntu/xrp_claude_bot/.env
EnvironmentFile=-/home/ubuntu/SYGNIF/.env
ExecStart=/home/ubuntu/SYGNIF/.venv/bin/python3 /home/ubuntu/SYGNIF/scripts/bybit_stream_monitor.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
