# systemd/cursor-agent-worker.service
[Unit]
Description=Cursor Cloud Agent private worker (SYGNIF repo)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/xrp_claude_bot
Environment=HOME=/home/ubuntu
# Same secrets & Cursor/finance vars as Docker, finance-agent, trade_overseer (CURSOR_API_KEY, BYBIT_*, etc.)
EnvironmentFile=-/home/ubuntu/xrp_claude_bot/.env
ExecStart=/home/ubuntu/.local/bin/agent worker start --worker-dir /home/ubuntu/xrp_claude_bot --management-addr 127.0.0.1:8093 --name sygnif-xrp-worker
Restart=on-failure
RestartSec=15
StandardOutput=append:/home/ubuntu/.local/share/cursor-agent/worker.log
StandardError=append:/home/ubuntu/.local/share/cursor-agent/worker.log

[Install]
WantedBy=multi-user.target
