# systemd/sygnif-dashboard-btc-terminal.service
[Unit]
Description=Sygnif BTC Terminal (prediction + training channel)
After=docker.service network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/SYGNIF
# Same order as compose ``x-sygnif-env-files``: secrets first, repo ``.env`` overrides.
EnvironmentFile=-/home/ubuntu/xrp_claude_bot/.env
EnvironmentFile=-/home/ubuntu/SYGNIF/.env
# Port 8888 — mutually exclusive with ``sygnif-dashboard-spot`` (same host port). Stop spot or set SYGNIF_BTC_TERMINAL_PORT in .env.
Environment=SYGNIF_BTC_TERMINAL_PORT=8888
ExecStart=/usr/bin/python3 /home/ubuntu/SYGNIF/dashboard_server_btc_terminal.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
