# systemd/sygnif-finetune-automation.service
[Unit]
Description=SYGNIF ML finetune loop (btc_predict_runner + channel_training; lock-shared with ruleprediction)
Documentation=file:///home/ubuntu/SYGNIF/scripts/sygnif_finetune_automation.sh
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
EnvironmentFile=-/home/ubuntu/SYGNIF/.env
ExecStart=/bin/bash /home/ubuntu/SYGNIF/scripts/sygnif_finetune_automation.sh --loop
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
