# ec2-snapshot/systemd/sygnif-brain-insights.service.d/disable-channel-write.conf
[Service]
# 2026-05-10: Option 2 path divergence fix — block SYGNIF-side processes from
# writing the legacy 5-source channel.json (which would clobber the rich
# 7-source content via the symlink). The rich writer is sygnif-bybit-nl-feed
# (sygnif-swarm tree) which intentionally does NOT have this env set.
Environment=SYGNIF_DISABLE_CHANNEL_WRITE=1
