# ec2-snapshot/systemd/sygnif-discovery.timer
[Unit]
Description=Run SYGNIF discovery pass every 30 minutes
Requires=sygnif-discovery.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
