# ec2-snapshot/systemd/sygnif-whale-alignment-audit.timer
[Unit]
Description=Run SYGNIF whale-alignment audit daily at 04:00 UTC
Requires=sygnif-whale-alignment-audit.service

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
