# ec2-snapshot/systemd/sygnif-perp-runner.service
[Unit]
Description=SYGNIF Perp Runner — scanner-driven perp executor with option-metric enrichment
After=network-online.target sygnif-bybit-daemon.service sygnif-microstructure-feed.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/sygnif-agent-mirror
EnvironmentFile=-/etc/sygnif/trader.env
EnvironmentFile=-/etc/sygnif/perp-runner.env
Environment=SYGNIF_PERP_DRY_RUN=1
ExecStart=/opt/sygnif/.venv/bin/python /home/ubuntu/sygnif-agent-mirror/perp_runner.py
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/ubuntu/.sygnif/perp-runner.log
StandardError=append:/home/ubuntu/.sygnif/perp-runner.log

[Install]
WantedBy=multi-user.target
