# ec2-snapshot/systemd/sygnif-intel-aggregator.service
[Unit]
Description=SYGNIF intel aggregator — pre-digests all intelligence into intel_summary.json
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=8

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/sygnif-services
Environment=PYTHONUNBUFFERED=1
Environment=SYGNIF_INTEL_POLL_S=30
EnvironmentFile=/etc/sygnif/chartinspect-keys.env
ExecStart=/opt/sygnif/.venv/bin/python -u /opt/sygnif-services/sygnif_intel_aggregator.py
StandardOutput=append:/var/log/sygnif/intel-aggregator.log
StandardError=append:/var/log/sygnif/intel-aggregator.log
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
