# ec2-snapshot/systemd/sygnif-trailing-manager.timer
[Unit]
Description=Run SYGNIF trailing-stop manager every 20 seconds
Requires=sygnif-trailing-manager.service

[Timer]
OnBootSec=60
OnUnitActiveSec=20
Persistent=true

[Install]
WantedBy=timers.target
