# ec2-snapshot/systemd/sygnif-btc-1h-refresh.timer
[Unit]
Description=Refresh btc_1h_ohlcv.json every 30 min (must finish before btc-predict-runner.timer fires)
Documentation=file:///home/ubuntu/SYGNIF/finance_agent/btc_specialist/scripts/refresh_1h_ohlcv.py

[Timer]
OnCalendar=*:0/30:00
AccuracySec=30s
Persistent=true
RandomizedDelaySec=20

[Install]
WantedBy=timers.target
