# ec2-snapshot/systemd/sygnif-telegram-relay.service
[Unit]
Description=SYGNIF DLP commentary -> Telegram relay
After=network-online.target sygnif-dlp.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/sygnif-agent-mirror
EnvironmentFile=-/etc/sygnif/trader.env
EnvironmentFile=-/etc/sygnif/telegram-relay.env
Environment=SYGNIF_TG_RELAY_INTERVAL=30
Environment=SYGNIF_TG_RELAY_SKIP_LABELS=routine
ExecStart=/opt/sygnif/.venv/bin/python /opt/sygnif-services/sygnif_telegram_relay.py
Restart=on-failure
RestartSec=15
StandardOutput=append:/home/ubuntu/.sygnif/tg-relay.log
StandardError=append:/home/ubuntu/.sygnif/tg-relay.log

[Install]
WantedBy=multi-user.target
