# ec2-snapshot/systemd/sygnif-decision-joiner.timer
[Unit]
Description=Run SYGNIF decision-joiner every hour
Requires=sygnif-decision-joiner.service

[Timer]
OnBootSec=300
OnUnitActiveSec=3600
Persistent=true

[Install]
WantedBy=timers.target
