# ec2-snapshot/systemd/sygnif-tron-signals.service
[Unit]
Description=SYGNIF Tron-side stablecoin mint tracker (USDT/USDC/USDD on TRON)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=8

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/sygnif-services
Environment=PYTHONUNBUFFERED=1
Environment=SYGNIF_TRON_MINT_USD=10000000
Environment=SYGNIF_TRON_POLL_S=300
EnvironmentFile=/etc/sygnif/tron-keys.env
ExecStart=/opt/sygnif/.venv/bin/python -u /opt/sygnif-services/sygnif_tron_signals.py
StandardOutput=append:/var/log/sygnif/tron-signals.log
StandardError=append:/var/log/sygnif/tron-signals.log
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
